# awlsim-server.service.in
[Unit]
Description=Awlsim soft-PLC core server
After=syslog.target
After=network.target

[Service]
Type=simple
TimeoutSec=60
Restart=on-failure
RestartSec=100ms
User=@USER@
Group=@GROUP@
Nice=-15

ExecStart=@PYTHON@ @PREFIX@/bin/awlsim-server -l localhost -4 -S -w @PROJECT@
ExecStartPost=-@PYTHON@ @PREFIX@/bin/awlsim-client -c localhost -r RUN

Environment=PYTHONPATH=@PYTHON_SITE@ PYTHONHASHSEED=0 PYTHONOPTIMIZE=2 PYTHONDONTWRITEBYTECODE=1 AWLSIM_CYTHON=1 AWLSIM_SCHED=realtime AWLSIM_PRIO=50 AWLSIM_AFFINITY=-1,-2,-3

[Install]
WantedBy=multi-user.target
