# awlsim-server.service.in
[Unit]
Description=Awlsim soft-PLC core server
After=syslog.target
After=network.target

[Service]
Type=simple
TimeoutSec=300
User=@USER@
Group=@GROUP@
Nice=-15

ExecStart=@PREFIX@/bin/awlsim-server -l any -w /home/@USER@/.core.awlpro
ExecStartPost=@PREFIX@/bin/awlsim-client -c localhost -r RUN

Environment=PYTHONPATH=@PYTHON_SITE@ AWLSIMCYTHON=1

[Install]
WantedBy=multi-user.target
